# fss.cron
#!/bin/sh

/usr/bin/smbclient -U fss -c "lcd /data/school/FSS;prompt;get Finger.pk" //sfb-f96f3ac27a8.borgia.com/fssbiowedge "$(cat /data/school/FSS/fssbiowedge.secret)"
/usr/bin/smbclient -U fss -c "lcd /data/school/FSS/POR;prompt;recurse;mget *" //sfb-f96f3ac27a8.borgia.com/StartUp "$(cat /data/school/FSS/fssbiowedge.secret)"
/usr/bin/smbclient -U fss -c "lcd /data/school/FSS/POR;prompt;recurse;mget SSH" //sfb-f96f3ac27a8.borgia.com/ApplicationData "$(cat /data/school/FSS/fssbiowedge.secret)"
/usr/bin/smbclient -U fss -c "lcd /data/school/FSS/POSS;prompt;recurse;mget *.zip" //windb.borgia.com/Backups$ "$(cat /data/school/FSS/fssbiowedge.secret)"
# Put Follett.pk on Library
